# s6-rc.d/ucsschool-id-connector-rest-api/run
#!/command/with-contenv /bin/bash

/usr/local/bin/uvicorn --host 0.0.0.0 --port 8911 ucsschool_id_connector.http_api:app
